# deploy/supervisord.render.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0

[program:ml_service]
directory=/app/ml_service
command=python -m uvicorn app:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:server]
directory=/app/server
command=node server.js
environment=ML_SERVICE_URL="http://127.0.0.1:8000"
autostart=true
autorestart=true
priority=20
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
